# files/etc/systemd/system/var-lib-systemd.mount
[Unit]
Documentation=man:fstab(5)
Before=local-fs.target
RequiresMountsFor=/sysimage/mutable

[Mount]
What=/sysimage/mutable/var/lib/systemd
Where=/var/lib/systemd
Type=none
Options=bind
